# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for staging-tdd.murrx.me"

start on net-device-up
stop on shutdown

respawn

setuid robin
chdir /home/robin/sites/staging-tdd.murrx.me/source

exec ../virtualenv/bin/gunicorn --bind unix:/tmpsock/staging-tdd.socket superlists.wsgi:application
